# config/etc/systemd/system/sosci_gunicorn.socket
[Unit]
Description=sosci gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
